# etc/systemd/system/docked.target
[Unit]
Description=Laptop is docked
Conflicts=undocked.target
