# etc/systemd/system/docked_at_office.target
[Unit]
Description=Laptop is docked at office
Conflicts=undocked.target
