# etc/systemd/system/undocked.target
[Unit]
Description=Laptop is undocked
Conflicts=docked.target
Conflicts=docked_at_office.target
